# osg_authentication/update_ca_certificates.cron
* * * * * root echo "Updating $(date)" >> /var/log/osg_ca_cron && yum update osg-ca-certs >> /var/log/osg_ca_cron
